# .config/systemd/user/gnome-session.target
[Unit]
Description=GNOME graphical session
BindsTo=graphical-session.target
