# .config/systemd/user/screen.service
[Install]
WantedBy=default.target

[Service]
ExecStart=/usr/bin/screen -Dm
Restart=always
