# .config/systemd/user/tmpfiles-clean.service
[Unit]
Description=Clean up temporary directories

[Service]
Type=oneshot
ExecStart=/bin/systemd-tmpfiles --clean %h/.config/systemd-tmpfiles.conf
