# .config/systemd/user/emacs.service.d/brett.conf
[Service]
ExecReload=/usr/bin/emacsclient --eval '(load-file (concat user-emacs-directory "init.el"))'
UMask=077
Environment=SSH_AUTH_SOCK=%t/keyring/ssh
Restart=on-failure
